# debian/onvif-recorder.service
[Unit]
Description=ONVIF Detection Recorder
Documentation=https://github.com/danielwoz/ubiquiti-protect-onvif-event-listener
After=network.target postgresql@14-protect.service unifi-protect.service
Wants=postgresql@14-protect.service

[Service]
Type=simple
User=root
# Shipped defaults (conffile; dpkg may prompt on upgrades if edited).
EnvironmentFile=-/etc/default/onvif-recorder
# User override (never shipped by the package; persists across upgrades).
# Variables set here win because systemd loads files in declared order.
EnvironmentFile=-/etc/default/onvif-recorder.local
ExecStart=/usr/bin/onvif-recorder $ONVIF_RECORDER_FLAGS
# `always` (rather than `on-failure`) lets the admin UI write a new config
# and exit cleanly to apply it; systemd brings the service back up.
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
